# supervisor/horizon.conf
[program:horizon]
process_name=%(program_name)s_%(process_num)02d
environment=CACHE_DRIVER=redis
command=/usr/local/bin/php /var/www/artisan horizon
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/tmp/horizon-out.log
stderr_logfile=/tmp/horizon-err.log
